# dns/supervisord.conf
[supervisord]
nodaemon=true

[program:named]
command=/usr/sbin/named -f
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/named.log
